# apt/supervisor/polipo.conf
[program:polipo]
command=polipo socksParentProxy=localhost:1080
autorestart=true
user=root
