# apt/supervisor/shadowsocks.conf
[program:shadowsocks]
command=/usr/local/bin/sslocal -c /etc/shadowsocks.json
autorestart=true
user=root
